# files/slideshow.service
[Unit]
Description=Slideshow

[Service]
ExecStart=/usr/local/bin/slideshow.sh

[Install]
WantedBy=multi-user.target
